# deployments/systemd/mpcium.service
[Unit]
Description=Mpcium MPC Node
Documentation=https://github.com/fystack/mpcium
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mpcium
Group=mpcium
WorkingDirectory=/opt/mpcium
ExecStart=/usr/local/bin/mpcium start --name ${MPCIUM_NODE_NAME} --config /etc/mpcium/config.yaml --password-file=%d/mpcium-db-password.cred
EnvironmentFile=/opt/mpcium/.env
SetCredentialEncrypted=mpcium-db-password.cred:DB_PASSWORD_BASE64_BLOB_DATA
SetCredentialEncrypted=mpcium-identity-password.cred:IDENTITY_PASSWORD_BASE64_BLOB_DATA
Restart=always
RestartSec=20
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mpcium
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ReadOnlyPaths=/etc/mpcium
ReadWritePaths=/opt/mpcium
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
